# .config/systemd/user/syncthing.service
[Unit]
Description=Syncthing
BindsTo=home.mount

[Service]
ExecStart=/usr/bin/syncthing

[Install]
WantedBy=home.mount
